# Vagrant_Ansible/Provisioning/files/gunicorn.ini
[program:webdev]
directory=/webapps/devops/devops-challenge/
command=/usr/bin/gunicorn app:app -b localhost:8000
user=vagrant
autostart=true
autorestart=true
stderr_logfile=/var/log/webdev/webdev.err.log
stdout_logfile=/var/log/webdev/webdev.out.log
